# support/systemd/lib/systemd/system/docker-volume-netshare.service
[Unit]
Description=Docker NFS, AWS EFS & Samba/CIFS Volume Plugin
Documentation=https://github.com/ContainX/docker-volume-netshare
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service
BindsTo=docker.service


[Service]
# Environment files: /etc/default/ takes precedence over /etc/sysconfig/
# The '-' prefix means the file is optional (no error if missing)
EnvironmentFile=-/etc/sysconfig/docker-volume-netshare
EnvironmentFile=-/etc/default/docker-volume-netshare
ExecStart=/usr/bin/docker-volume-netshare $DKV_NETSHARE_OPTS
Restart=on-failure
RestartSec=5s
StandardOutput=journal

[Install]
WantedBy=multi-user.target
